# conf/etc/s6/services/03-slapd/run
#!/usr/bin/with-contenv bash
ulimit -n 1024
service-handler "OPENLDAP" "/usr/sbin/slapd -h 'ldap://' -u ldap -g ldap -d $LOG_LEVEL" "/tmp/state/10-openldap-init"
